# php54-supervisor/build/files/etc/supervisor/sample.conf
[program:sample]
command=bash -c 'echo "THIS IS SAMPLE SUPERVISOR JOB" && sleep 120'
process_name=%(program_name)s_%(process_num)02d
numprocs=1
exitcodes=0
autostart=true
autorestart=true
user=project
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

